# install_files/etc/cron.hourly/rsync_www_hourly
#!/bin/bash

RSYNC_USER="[RSYNC_USER]"
BACKUP_PRIVATEIP="[OTHER_PRIVATEIP]"

LOCAL_WWW="/var/www/"
REMOTE_WWW="~/www/"
RSYNC_BIN="$(which rsync)"

$RSYNC_BIN -azP $LOCAL_WWW $RSYNC_USER@$BACKUP_PRIVATEIP:$REMOTE_WWW

